# syslog_service.service
[Unit]
Description=Syslog Retention and SIEM Service
After=network.target rsyslog.service
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/syslog-retention-service
ExecStart=/opt/syslog-retention-service/.venv/bin/python main.py
Restart=on-failure
RestartSec=30
# Prevent rapid restart loops from consuming all CPU
StartLimitIntervalSec=120
StartLimitBurst=3
# Memory and CPU guards to prevent runaway resource use
MemoryMax=512M
CPUQuota=80%
StandardOutput=append:/var/log/syslog-siem/service.log
StandardError=append:/var/log/syslog-siem/service_err.log

[Install]
WantedBy=multi-user.target
